# supervisor/conf.d/rtl_915.conf
[program:rtl_915]
priority=10
user=docker
group=nogroup
command=/docker/rtl_433-to-influx.py /docker/config_915a.yaml
stopsignal=INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
