# supervisor/conf.d/rtl_433.conf
[program:rtl_433]
priority=10
user=docker
command=/docker/rtl_433-to-influx.py
group=nogroup
stopsignal=INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
